# phy_parser.service
[Unit]
Description=PHY Information

[Service]
Type=simple
ExecStart=/etc/inventec_svt_deployment/phy_parser

[Install]
WantedBy=multi-user.target
