# ping.service
[Unit]
Description=PING services

[Service]
Type=simple
ExecStart=/etc/inventec_svt_deployment/ping

[Install]
WantedBy=multi-user.target
