# ru_acc_parser.service
[Unit]
Description=RU & ACC Parser

[Service]
Type=simple
ExecStart=/etc/inventec_svt_deployment/ru_acc_parser

[Install]
WantedBy=multi-user.target
